# com.claude-code-releases.weekly-update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>HOME</key>
		<string>/Users/tairyu</string>
		<key>PATH</key>
		<string>/usr/local/bin:/usr/bin:/bin:/Users/tairyu/.local/bin</string>
	</dict>
	<key>Label</key>
	<string>com.claude-code-releases.weekly-update</string>
	<key>Program</key>
	<string>/bin/bash</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Users/tairyu/program/claude-code-releases/scripts/weekly-update.sh</string>
	</array>
	<key>StandardErrorPath</key>
	<string>/Users/tairyu/program/claude-code-releases/logs/weekly-launchd-stderr.log</string>
	<key>StandardOutPath</key>
	<string>/Users/tairyu/program/claude-code-releases/logs/weekly-launchd-stdout.log</string>
	<key>StartCalendarInterval</key>
	<array>
		<dict>
			<key>Hour</key>
			<integer>14</integer>
			<key>Minute</key>
			<integer>30</integer>
			<key>Weekday</key>
			<integer>4</integer>
		</dict>
	</array>
	<key>WorkingDirectory</key>
	<string>/Users/tairyu/program/claude-code-releases</string>
</dict>
</plist>
